# configs/supervisor.ini
[program:flask]
directory = /main
command = gunicorn
autostart = true
autorestart = true
startretries = 3


[program:huey]
directory = /main
command = huey_consumer app_flask.huey.run
autostart = true
autorestart = true
startretries = 3
